# chef/templates/default/init.conf.erb
env DAEMON=/opt/<%= @projectName %>/<%= @projectName %>-<%= @projectVersion %>/bin/<%= @projectName %>
env PID=/opt/<%= @projectName %>/<%= @projectName %>-<%= @projectVersion %>/RUNNING_PID

start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [!2345]

pre-start script
  # main Play process will not be run if $DAEMON does not exist
  test -f $DAEMON || { stop ; exit 0; }
end script

exec $DAEMON -mem 512 -J-server
